# systemd/system/strelaysrv.service
[Unit]
Description=Relay service for syncthing

[Service]
ExecStart=/usr/bin/strelaysrv
Restart=on-failure

[Install]
WantedBy=default.target
